# kbp_puppetmaster/files/puppetstoredconfigcleanhenker.cron
#!/bin/sh

exec puppetstoredconfigcleanhenker /etc/mysql/debian.cnf puppet
