# kbp_webalizer/templates/cron.daily/webalizer
#!/bin/sh

/usr/local/bin/webalizer-multi

# Exit with the last exit code
exit $?
