# kbp_puppetmaster/templates/puppetstoredconfigcleanhenker.cron
#!/bin/sh

exec puppetstoredconfigcleanhenker /etc/mysql/debian.cnf <%= name %><% if configfile %> <%= configfile %><% end %>
